# authelia.service
[Unit]
Description=Authelia authentication and authorization server
After=network.target

[Service]
ExecStart=/usr/bin/authelia --config /etc/authelia/configuration.yml
SyslogIdentifier=authelia

[Install]
WantedBy=multi-user.target
