# ci/root/etc/services.d/buildkite/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid buildkite ssh-env-config.sh /usr/local/bin/buildkite-agent start