# ci/root/etc/services.d/dockerperms/run
#!/usr/bin/with-contenv bash

while [ ! -S "/run/docker.sock" ];
do
  sleep 1;
done
chown root:abc /run/docker.sock
tail -f /dev/null